# images/full-node/crontab
0 * * * * chompers chompchain --generate
