# images/full-node/crontab
TZ = America/New_York
* * * * * chompers chompchain --generate
