# images/full-node/crontab
* * * * * chompers bash -c 'source ~/.bashrc && chompchain --generate'
