# images/full-node/crontab
* * * * * chompers bash -c 'source /home/chompers/.bashrc && /usr/local/bin/chompchain --generate'
