# images/demo-cron/crontab
TZ = America/New_York
* * * * * root echo "cron -- $(date)" >> /var/log/cron.log 2>&1
